# etc/systemd/system/rpm-ostreed-automatic.service.d/99-override-on-battery.conf
[Unit]
# Only upgrade system if we're connected to AC power
ConditionACPower=true
